# crontab
20 * * * * /home/pi/monip/monip.sh > /tmp/monip.log 2>&1
